# oam/src/eddie.init.src
#!/bin/sh
#
# Startup script for %PRODUCT_NAME%
#
# chkconfig: - 85 15
# description: %PRODUCT_NAME%
#
# 

PATH=$PATH:/usr/local/lib/erlang/bin:/usr/local/sbin

# See how we were called.
case "$1" in
  start)
	echo -n "Starting %PRODUCT_NAME%: "
	eddie daemon
	echo
	;;
  stop)
	echo -n "Suspending %PRODUCT_NAME%: "
	eddie suspend
	echo
	;;
  reconfig)
	echo -n "Reconfiguring %PRODUCT_NAME%: "
	eddie config
	echo
	;;
  *)
	echo "Usage: $0 {start|stop|reconfig}"
	exit 1
esac

exit 0
